# deploy/com.tangleclaw.ttyd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tangleclaw.ttyd</string>
    <key>ProgramArguments</key>
    <array>
        <string>__TTYD_PATH__</string>
        <string>--writable</string>
        <string>--url-arg</string>
        <string>--port</string>
        <string>3100</string>
        <string>__REPO_DIR__/deploy/ttyd-attach.sh</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/tangleclaw-ttyd.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/tangleclaw-ttyd.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__LAUNCHD_PATH__</string>
    </dict>
</dict>
</plist>
