# deploy/com.tangleclaw.ttyd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tangleclaw.ttyd</string>
    <key>ProgramArguments</key>
    <array>
        <!-- #397 bug 2 (+ durability): the launchd PROGRAM must be a non-TCC
             system binary, never a script under the repo. When the repo lives in
             a TCC-protected tree (~/Documents — the documented default clone
             location), a launchd-spawned process cannot read-to-exec a script
             there without Full Disk Access → exit 126 and ttyd never starts (the
             server agent survives only because its program is the node binary,
             outside ~/Documents). So we run /bin/bash and carry the self-heal
             INLINE FROM ARGV (no protected-file read): unlink a stale Unix socket
             (caddy mode; a no-op when TTYD_SOCKET is empty, i.e. direct mode),
             then exec the ttyd binary. Both /bin/bash and ttyd are non-TCC.
             NOTE (#500): ttyd-attach.sh (__TTYD_ATTACH__ below) is installed
             OUTSIDE the repo, at ~/.tangleclaw/deploy/. It is a ttyd *arg* that
             ttyd opens per client-connect — and ttyd is DENIED Full Disk Access,
             so a repo-resident copy under ~/Documents freezes ttyd in open()
             (all sessions black-screen after a ttyd restart, 2026-07-06). Same
             TCC family as the launchd-program path above, a *different* syscall
             site. install.sh / ingress-cutover.js sync the script to the non-TCC
             path (lib/ttyd-attach.js) and fill __TTYD_ATTACH__ with it. -->
        <string>/bin/bash</string>
        <string>-c</string>
        <string>if [ -n "${TTYD_SOCKET:-}" ]; then rm -f "$TTYD_SOCKET"; fi; exec "$@"</string>
        <string>ttyd-launch</string>
        <string>__TTYD_PATH__</string>
        <string>--writable</string>
        <string>--url-arg</string>
        <!-- AUTH-1 (#395): bind args are filled per ingress mode — direct mode
             (install.sh) → `--port 3100`; caddy mode (ingress-cutover.js) →
             `--interface <unix-socket>` so ttyd is reachable only via the proxy
             chain behind Caddy. -->
        <string>__TTYD_BIND_KEY__</string>
        <string>__TTYD_BIND_VAL__</string>
        <!-- Raise xterm.js scrollback from its 1000-line default (#322/#290):
             gives long output room within a live session, and is large enough
             to hold the history ttyd-attach.sh replays on reconnect/restart. -->
        <string>--client-option</string>
        <string>scrollback=10000</string>
        <string>__TTYD_ATTACH__</string>
    </array>
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>
    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/dev/null</string>
    <key>StandardErrorPath</key>
    <string>/dev/null</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__LAUNCHD_PATH__</string>
        <!-- #397 bug 2: the ttyd Unix socket the inline launcher unlinks before
             bind (caddy mode). Filled with the socket path by ingress-cutover.js;
             empty in direct mode (install.sh) where ttyd binds a TCP port. -->
        <key>TTYD_SOCKET</key>
        <string>__TTYD_SOCKET__</string>
    </dict>
</dict>
</plist>
