# deploy/com.tangleclaw.ttyd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tangleclaw.ttyd</string>
    <key>ProgramArguments</key>
    <array>
        <string>__TTYD_PATH__</string>
        <string>--writable</string>
        <string>--url-arg</string>
        <!-- AUTH-1 (#395): bind args are filled per ingress mode — direct mode
             (install.sh) → `--port 3100`; caddy mode (ingress-cutover.js) →
             `--interface <unix-socket>` so ttyd is reachable only via the proxy
             chain behind Caddy. -->
        <string>__TTYD_BIND_KEY__</string>
        <string>__TTYD_BIND_VAL__</string>
        <!-- Raise xterm.js scrollback from its 1000-line default (#322/#290):
             gives long output room within a live session, and is large enough
             to hold the history ttyd-attach.sh replays on reconnect/restart. -->
        <string>--client-option</string>
        <string>scrollback=10000</string>
        <string>__REPO_DIR__/deploy/ttyd-attach.sh</string>
    </array>
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>
    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/dev/null</string>
    <key>StandardErrorPath</key>
    <string>/dev/null</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__LAUNCHD_PATH__</string>
    </dict>
</dict>
</plist>
